# ups-battery-watchdog.timer
[Unit]
Description=Run UPS battery watchdog every minute
Documentation=https://github.com/moijafcor/ups-battery-watchdog

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=15s

[Install]
WantedBy=timers.target
